# devops/raptor-launcher.service
[Unit]
Description=Raptor Launcher - Erlang Service Orchestrator
After=network.target postgresql.service docker.service
Wants=postgresql.service docker.service

[Service]
Type=simple
User=matkat
Group=matkat
WorkingDirectory=/home/matkat/raptor-launcher/launcher-app/raptor_launcher
Environment="HOME=/home/matkat"

# Load environment variables from files
# Note: systemd EnvironmentFile only supports KEY=VALUE pairs
# No bash syntax (export, if, $(), etc.)
EnvironmentFile=/home/matkat/launcher-app/devops/env/raptor.systemd.env
EnvironmentFile=-/home/matkat/launcher-app/devops/env/.systemd.env

# Start command - using Erlang release (not interactive shell!)
ExecStart=/home/matkat/launcher-app/raptor_launcher foreground

# Restart policy
Restart=on-failure
RestartSec=10
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raptor-launcher

# Security settings (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
