# devops/raptor-launcher.service
[Unit]
Description=Raptor Launcher - Erlang Service Orchestrator
After=network.target postgresql.service docker.service
Wants=postgresql.service docker.service

[Service]
Type=simple
User=matkat
Group=matkat
WorkingDirectory=/home/matkat/launcher-app
Environment="HOME=/home/matkat"

# Start command - using wrapper script
EnvironmentFile=/home/matkat/launcher-app/devops/env/raptor.systemd.env
EnvironmentFile=-/home/matkat/launcher-app/devops/env/.systemd.env
ExecStart=/home/matkat/launcher-app/bin/start_app.sh

# Restart policy
Restart=on-failure
RestartSec=10
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raptor-launcher

# Security settings (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
